# config/agent.cron
* * * * *    /bin/touch /opt/uravo/run/crond.running
22 */6 * * * /opt/uravo/bin/update_uravo.pl
*/5 * * * *  /opt/uravo/bin/agent.pl
