# config/outpost.cron
* * * * * /opt/uravo/bin/outpost.pl &
*/5 * * * * /opt/uravo/bin/remote_agent.pl

